# platform/broadcom/sonic-platform-modules-nexthop/common/etc/systemd/system/determine-reboot-cause.service.d/override.conf
[Unit]
Requires=pddf-platform-init.service
After=pddf-platform-init.service

[Service]
# Skip if already ran this boot, while keeping the service in "active" state.
# Do NOT use ConditionPathExists or ExecCondition here. They leave the service
# "inactive" (skipped), which causes dependents like process-reboot-cause to fail.
ExecStart=
ExecStart=/bin/sh -c 'test -f /run/determine-reboot-cause-done || /usr/local/bin/determine-reboot-cause'
ExecStartPost=/usr/bin/touch /run/determine-reboot-cause-done
